# deploy/systemd/cloudflared-openqlow.service
[Unit]
Description=Cloudflare named tunnel for OPENQLOW
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
User=openqlow
Group=openqlow
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/cloudflared/config.yml run
Restart=always
RestartSec=10
MemoryMax=128M
CPUQuota=20%
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadOnlyPaths=/etc/cloudflared
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
