# deploy/systemd/openqlow-monitor.service
[Unit]
Description=OPENQLOW health monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/openqlow
EnvironmentFile=/etc/openqlow/openqlow.env
Environment=OPENQLOW_MONITOR_SYSTEMD=true
Environment=OPENQLOW_MONITOR_SELF_HEAL=true
Environment=OPENQLOW_MONITOR_SERVICES=openqlow-webhook.service
ExecStart=/usr/bin/npm run monitor
MemoryMax=256M
CPUQuota=25%
TasksMax=64
OOMPolicy=stop
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/openqlow /opt/obsidian-vault
