# deploy/systemd/openqlow-crm-daily-report.service
[Unit]
Description=OPENQLOW CRM daily report generation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=openqlow
Group=openqlow
WorkingDirectory=/opt/openqlow
EnvironmentFile=/etc/openqlow/openqlow.env
ExecStart=/usr/bin/node dist/crm/crm_cli.js daily-report
MemoryMax=256M
CPUQuota=30%
TasksMax=128
OOMPolicy=stop
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/openqlow /home/flatup/openqlow-data
StandardOutput=journal
StandardError=journal
