# deploy/systemd/openqlow-morning.service
[Unit]
Description=OPENQLOW morning briefing push (Jin LINE)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=openqlow
Group=openqlow
WorkingDirectory=/opt/openqlow
EnvironmentFile=/etc/openqlow/openqlow.env
ExecStart=/usr/bin/node dist/scheduler/morning_briefing.js
MemoryMax=192M
CPUQuota=30%
TasksMax=64
OOMPolicy=stop
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/openqlow
