# deploy/systemd/openqlow-daily.service
[Unit]
Description=OPENQLOW daily idea generation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=openqlow
Group=openqlow
WorkingDirectory=/opt/openqlow
EnvironmentFile=/etc/openqlow/openqlow.env
ExecStart=/usr/bin/node dist/scheduler/daily.js
MemoryMax=512M
CPUQuota=60%
TasksMax=128
OOMPolicy=stop
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/openqlow /opt/obsidian-vault
