# deploy/systemd/openqlow-daily.timer
[Unit]
Description=Run OPENQLOW daily at 05:15 JST

[Timer]
OnCalendar=*-*-* 05:15:00 Asia/Tokyo
Persistent=true
Unit=openqlow-daily.service

[Install]
WantedBy=timers.target
